# debian/config/supervisord.conf
# =============================================================================
# Supervisor Configuration for KAZMA Desktop
# =============================================================================
# Manages XRDP and related services
# =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
logfile_maxbytes=10MB
logfile_backups=1
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# =============================================================================
# XRDP Session Manager
# =============================================================================
[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman --nodaemon
autostart=true
autorestart=true
priority=100
stdout_logfile=/var/log/supervisor/xrdp-sesman.log
stderr_logfile=/var/log/supervisor/xrdp-sesman-error.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB

# =============================================================================
# XRDP Server
# =============================================================================
[program:xrdp]
command=/usr/sbin/xrdp --nodaemon
autostart=true
autorestart=true
priority=200
startsecs=5
stdout_logfile=/var/log/supervisor/xrdp.log
stderr_logfile=/var/log/supervisor/xrdp-error.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB

